# 85-sleepi-utils.rules
KERNEL=="rtc*", RUN+="/lib/udev/sleepi-utils.sh $name"
